# powersentry.service
[Unit]
Description=PowerSentry — Virtual BMS Watchtower for Inverter + Lead-Acid Battery
Documentation=https://github.com/think777/powersentry
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=user
WorkingDirectory=/home/user/powersentry
ExecStart=/home/user/powersentry/.venv/bin/python powersentry.py
Restart=on-failure
RestartSec=30


# Logging — stdout/stderr go to journalctl
StandardOutput=journal
StandardError=journal
SyslogIdentifier=powersentry

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome=read-only — removed, conflicts with WorkingDirectory under /home/
ReadWritePaths=/home/user/powersentry
PrivateTmp=true

[Install]
WantedBy=multi-user.target
